# container/services/resty/run
#!/bin/sh
exec 2>&1

exec /usr/local/openresty/bin/openresty -g 'daemon off;' -c /container/config/resty/nginx.conf
